# babybuddy/root/etc/services.d/babybuddy/run
#!/usr/bin/with-contenv bashio
cd /app/babybuddy || (echo "/app/babybuddy does not exist, exiting..." & exit)

export \
    DJANGO_SETTINGS_MODULE="babybuddy.settings.base" \
    ALLOWED_HOSTS="${ALLOWED_HOSTS:-*}" \
    TIME_ZONE="${TZ:-UTC}" \
    DEBUG="${DEBUG:-False}" \
    SECRET_KEY="${SECRET_KEY:-`cat /config/.secretkey`}"

export CSRF_TRUSTED_ORIGINS=$(bashio::config "CSRF_TRUSTED_ORIGINS")

exec \
    s6-setuidgid root gunicorn babybuddy.wsgi -b :8000 --log-level=info \
        --worker-tmp-dir=/dev/shm --log-file=- \
        --workers=2 --threads=4 --worker-class=gthread \
