# babybuddy/root/etc/services.d/babybuddy/run
#!/command/with-contenv bashio
#

declare -a options

options+=(-b :8000)
options+=(--log-level="$(bashio::config 'log_level')")
options+=(--worker-tmp-dir=/dev/shm)
options+=(--log-file=-)
options+=(--workers=2)
options+=(--threads=4)
options+=(--worker-class=gthread)

cd /app/babybuddy || bashio.exit.nok 'Failed cd, /app/babybuddy does not exist, exiting...'

export \
    DJANGO_SETTINGS_MODULE="babybuddy.settings.base" \
    ALLOWED_HOSTS="${ALLOWED_HOSTS:-*}" \
    TIME_ZONE="${TZ:-UTC}" \
    DEBUG="${DEBUG:-False}" \
    ENABLE_HOME_ASSISTANT_SUPPORT=true \
    SECRET_KEY="${SECRET_KEY:-$(cat /config/.secretkey)}"

export \
    CSRF_TRUSTED_ORIGINS=$(bashio::config "CSRF_TRUSTED_ORIGINS")

if bashio::config.has_value "INGRESS_USER"; then
    bashio::log.info 'Adding config for Ingress User Auth'
    export \
        REVERSE_PROXY_AUTH=true \
        PROXY_HEADER=HTTP_X_REMOTE_USER
fi

exec \
    s6-setuidgid root gunicorn babybuddy.wsgi "${options[@]}"

