# babybuddy/rootfs/etc/services.d/babybuddy/run
#!/command/with-contenv bashio
#

declare -a options

options+=(-b :8000)
options+=(--log-level=info)
options+=(--worker-tmp-dir=/dev/shm)
options+=(--log-file=-)
options+=(--workers=2)
options+=(--threads=4)
options+=(--worker-class=gthread)

cd /app/babybuddy || bashio.exit.nok 'Failed cd, /app/babybuddy does not exist, exiting...'

export \
    DJANGO_SETTINGS_MODULE="babybuddy.settings.base" \
    ALLOWED_HOSTS="${ALLOWED_HOSTS:-*}" \
    TIME_ZONE="${TZ:-UTC}" \
    DEBUG="${DEBUG:-False}" \
    SECRET_KEY="${SECRET_KEY:-`cat /config/.secretkey`}"

export \
    CSRF_TRUSTED_ORIGINS=$(bashio::config "CSRF_TRUSTED_ORIGINS")

exec \
    s6-setuidgid root gunicorn babybuddy.wsgi "${options[@]}"
