# babybuddy/rootfs/etc/services.d/nginx/run
#!/command/with-contenv bashio

bashio::log.info "NGINX waiting on babybuddy..."
bashio::net.wait_for 8000

bashio::log.info "Starting NGinx..."
exec nginx
