# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
bashio::log.info "NgINX waiting on babybuddy at port 8889"
bashio::net.wait_for 8889

bashio::log.info "Starting NGinx..."
exec nginx
