# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# wait for babybuddy to be running on port 8000
bashio::net.wait_for 8000

bashio::log.info "Starting NGinx..."
exec nginx
